# scripts/setup/com.user.claude-stats-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.user.claude-stats-sync</string>

	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>source ~/.bashrc 2>/dev/null || source ~/.zshrc 2>/dev/null; exec ~/dotfiles/scripts/development/sync-claude-stats-to-icloud.ts</string>
	</array>

	<key>StartInterval</key>
	<integer>3600</integer>

	<key>RunAtLoad</key>
	<true/>

	<key>StandardOutPath</key>
	<string>/var/log/claude-stats-sync.log</string>

	<key>StandardErrorPath</key>
	<string>/var/log/claude-stats-sync.log</string>

	<key>EnvironmentVariables</key>
	<dict>
		<key>HOME</key>
		<string>$HOME</string>
		<key>PATH</key>
		<string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
	</dict>

	<key>Nice</key>
	<integer>10</integer>

	<key>ProcessType</key>
	<string>Background</string>
</dict>
</plist>
